# cloudera-cdh-hivemetastore/supervisord.conf
[supervisord]
nodaemon=true

[program:hive-metastore]
directory=/data
command=hive --service metastore
stdout_logfile=/var/log/hive/hive-metastore.log
redirect_stderr=true
user=hive
