# cloudera-cdh-datanode/supervisord.conf
[supervisord]
nodaemon=true

[program:nodemanager]
command=yarn nodemanager
redirect_stderr=true
stdout_logfile=/var/log/hadoop-yarn/yarn-nodemanager.log
user=yarn
