# cloudera-cdh-prestodb/supervisord.conf
[supervisord]
nodaemon=true

[program:presto]
command=/presto/presto-server-PRESTO_VER/bin/launcher run
stdout_logfile=/var/log/presto/presto.log
redirect_stderr=true
user=presto
